# scripts/tokengauge-daemon.service
[Unit]
Description=TokenGauge daemon (long-lived fetcher + socket server)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
# The systemd user manager's PATH omits ~/.local/bin, where ccusage and the
# tokengauge binaries usually live; without this the daemon's ccusage cost
# fetches fail and rows show usage without cost enrichment (usage limits are
# fetched natively and are unaffected).
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.local/bin/tokengauge-waybar --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
